# salt/configure/templates/unit.service.j2
[Unit]
Description=Soundfleet {{ unit }} service

[Service]
ExecStart={{ venv_dir }}/bin/{{ unit }}d
Restart=always
RestartSec=180

[Install]
WantedBy=default.target
